# roles/postgresql/templates/postgresql.service.j2
[Service]
Environment=PGPORT={{ postgresql_connection.port }}
Environment=PGDATA={{ postgresql.data_directory }}
